# appium.conf
[program:appium]
command=/usr/bin/wait-for-it localhost:8100 -t 60 -- /usr/bin/appium
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
